# scripts/crontab
*/5 * * * * /home/yousef/.dotfiles/scripts/battery-notify.sh
